# docs/guides/installation/configuration-files/supervisor/baserow-frontend.conf
[supervisord]
nodaemon = true

[program:nuxt]
environment =
    PRIVATE_BACKEND_URL='http://localhost:8000',
    PUBLIC_WEB_FRONTEND_URL='https://FRONTEND_DOMAIN',
    PUBLIC_BACKEND_URL='https://BACKEND_DOMAIN'
directory = /baserow/web-frontend
command = sh -c './node_modules/.bin/nuxt start --hostname 127.0.0.1 --config-file ./config/nuxt.config.local.js'
stdout_logfile=/var/log/baserow/frontend.log
stderr_logfile=/var/log/baserow/frontend.error
