# docs/guides/installation/configuration-files/supervisor/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/var/log/baserow/nginx.out
stderr_logfile=/var/log/baserow/nginx.error
autostart=true
autorestart=true
numprocs=1